# config/systemd/celestia-light.service.d/override.conf
[Service]
CPUWeight=50
MemoryLimit=2G
IOWeight=50